# conf/salt/project/worker/celery.conf
[program:{{ pillar['project_name'] }}-{{ name }}]
command={{ virtualenv_root }}/bin/celery -A {{ pillar['project_name'] }} {{ command }} {{ flags }}
user={{ pillar['project_name'] }}
directory={{ directory }}
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile={{ log_dir }}/{{ name }}.log
redirect_stderr=true
stderr_logfile={{ log_dir }}/{{ name }}.error.log
startsecs=1
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs=60
environment=DJANGO_SETTINGS_MODULE="{{ settings }}",
    {%- for key, value in pillar.get('secrets', {}).items() + pillar.get('env', {}).items() -%}
        {{ key }}="{{ value }}"{%- if not loop.last -%},{%- endif -%}
    {%- endfor -%}